# tools/launchd/com.nous.ccd-orphan-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nous.ccd-orphan-watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/madia/Documents/Projects/Nous AGaaS/Nous/tools/ccd_orphan_watchdog.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/madia</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>CCD_WATCHDOG_THRESHOLD_HOURS</key>
        <string>2</string>
    </dict>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>/Users/madia/nous-agaas/logs/ccd-orphan-watchdog.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/madia/nous-agaas/logs/ccd-orphan-watchdog.err.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
