# tools/launchd/com.nous.library-graph.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.nous.library-graph</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/madia</string>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>NOUS_WIKI</key>
		<string>/Users/madia/Documents/Projects/Nous AGaaS/Nous</string>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-lc</string>
		<string>cd "/Users/madia/Documents/Projects/Nous AGaaS/Nous" &amp;&amp; /opt/homebrew/bin/python3 tools/library_drain_queue.py --wiki "/Users/madia/Documents/Projects/Nous AGaaS/Nous" --prefer local --summary-json --max-files 25 --git-writeback --git-push-remotes vps,github</string>
	</array>
	<key>StandardOutPath</key>
	<string>/Users/madia/Documents/Projects/Nous AGaaS/Nous/logs/library-drain-launchd.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/madia/Documents/Projects/Nous AGaaS/Nous/logs/library-drain-launchd.err.log</string>
	<key>StartInterval</key>
	<integer>180</integer>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
