# launchd/com.nous.daily-skill-evals.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.nous.daily-skill-evals</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/madia</string>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>WIKI</key>
		<string>/Users/madia/nous-agaas/wiki</string>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/madia/nous-agaas/wiki/tools/daily-skill-evals.sh</string>
	</array>
	<key>StandardErrorPath</key>
	<string>/Users/madia/nous-agaas/logs/daily-skill-evals.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/madia/nous-agaas/logs/daily-skill-evals.out.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>4</integer>
		<key>Minute</key>
		<integer>15</integer>
	</dict>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
