# launchd/com.nous.satory-camera-doctor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nous.satory-camera-doctor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>-m</string>
        <string>agents.camera_doctor.main</string>
        <string>--tenant</string>
        <string>satory</string>
        <string>--dry-run</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/madia/nous-agaas/wiki</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/madia</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONPATH</key>
        <string>/Users/madia/nous-agaas/wiki</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/madia/nous-agaas/logs/satory-camera-doctor/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/madia/nous-agaas/logs/satory-camera-doctor/launchd.err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
